# scripts/keepalive.service
[Unit]
Description=Agents Workspace — Keep-alive para backends remotos
After=network.target

[Service]
Type=simple
WorkingDirectory=/path/to/agents-workspace
ExecStart=/path/to/agents-workspace/venv/bin/python scripts/keepalive.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
Environment=KEEPALIVE_INTERVAL=300

[Install]
WantedBy=multi-user.target
